# crates/spn/assets/systemd/spn-daemon.service
[Unit]
Description=SuperNovae Package Manager Daemon
Documentation=https://github.com/supernovae-st/supernovae-cli
After=network.target

[Service]
# Run as user service (not root)
Type=simple

# Binary location (will be replaced during install)
ExecStart=${SPN_BIN} daemon start --foreground

# Restart on failure
Restart=on-failure
RestartSec=5

# Logging
StandardOutput=append:%h/.spn/daemon.log
StandardError=append:%h/.spn/daemon.err

# Environment
Environment=HOME=%h
Environment=PATH=/usr/local/bin:/usr/bin:/bin:%h/.cargo/bin

# Working directory
WorkingDirectory=%h

# Resource limits
Nice=5
MemoryHigh=256M
MemoryMax=512M

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.spn
PrivateTmp=true

[Install]
WantedBy=default.target
