# assets/launchd/com.supernovae.spn-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Service identifier -->
    <key>Label</key>
    <string>com.supernovae.spn-daemon</string>

    <!-- Program to run -->
    <key>ProgramArguments</key>
    <array>
        <string>${SPN_BIN}</string>
        <string>daemon</string>
        <string>start</string>
        <string>--foreground</string>
    </array>

    <!-- Start at login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart on failure (but not on clean exit) -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- Throttle restarts (10 second minimum between launches) -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>${HOME}/.spn/daemon.log</string>
    <key>StandardErrorPath</key>
    <string>${HOME}/.spn/daemon.err</string>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>${HOME}</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:${HOME}/.cargo/bin</string>
    </dict>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>${HOME}</string>

    <!-- Nice priority (0 = normal, positive = lower priority) -->
    <key>Nice</key>
    <integer>5</integer>

    <!-- Process type (Background = lower priority than apps) -->
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
